# crontab/hw.sh
#!/bin/bash
echo "Hellos World">> /home/ec2-user/tp/script/hello_world.log
sleep 30
echo "Hello World" >> /home/ec2-user/tp/script/hello_world.log

